# root/etc/s6-overlay/s6-rc.d/init-qbittorrent-monitor-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "**** Configuring QBittorrent Error Monitor ****"

# Create directories if they don't exist
mkdir -p /config/logs

# Copy default config if it doesn't exist
if [[ ! -f /config/config.yaml ]]; then
    echo "**** Copying default config ****"
    cp /defaults/config.yaml /config/config.yaml
fi

# Create logs directory in config
mkdir -p /config/logs

# Set permissions
lsiown abc:abc /config
lsiown abc:abc /config/logs
lsiown -R abc:abc /app

echo "**** Configuration complete ****"
