# root/etc/s6-overlay/s6-rc.d/svc-qbittorrent-monitor/run
#!/command/with-contenv bash

cd /app
exec s6-setuidgid abc python3 /app/qbittorrent-monitor.py
